# Fusuma/fusuma.service
[Unit]
Description=Fusuma touch pad gestures

[Service]
ExecStart=~/Fusuma.sh

[Install]
WantedBy=multi-user.target
